# prepper.service
[Unit]
Description=Prepper D-Day service
Requires=Network.target
After=Network.target

[Service]
ExecStart=/usr/bin/prepper
